# deploy/enzyme-frontend.service
[Unit]
Description=Enzyme Search Engine Frontend
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/bingxing2/home/scx9kyg/xqc/enzyme-search-engine
ExecStart=/usr/bin/python3 -m http.server 8020 --bind 127.0.0.1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
